# mediamtx/etc/systemd/system/mediamtx.movie.service
[Unit]
Description=MediaMTX webcam

[Service]
Type=simple
ExecStart=/bin/bash /data/mediamtx/movie.sh

[Install]
WantedBy=multi-user.target