# home-assistant/etc/systemd/system/docker-wake-on-lan.service
[Unit]
Description=Listen for wake-on-lan from docker via pipe

[Service]
Restart=always
ExecStart=/usr/bin/wakeonlan -f /data/home-assistant/config/wolpipe

[Install]
WantedBy=multi-user.target
